# deploy/systemd/caddy.service.d/admin-token.conf
[Service]
EnvironmentFile=-/etc/caddy/secrets.env
