# deploy/token-dashboard.service
[Unit]
Description=DelegateAgent Token Dashboard (per-group cost analytics)
After=network.target docker.service delegate-agent.service
Requires=docker.service
Wants=delegate-agent.service

[Service]
# `docker compose up` runs in the foreground so systemd can supervise it.
# Build + image pull happen in deploy/post-deploy.sh on each deploy, not here —
# that keeps build output visible in deploy logs and avoids systemd's lack of
# shell support in ExecStartPre. If you `systemctl restart`, the existing
# container is gracefully recreated. Crash recovery on boot is double-covered:
# this Restart=always plus `restart: unless-stopped` in docker-compose.yml.
Type=simple
User=root
WorkingDirectory=/opt/delegate-agent/deploy/token-dashboard
EnvironmentFile=-/opt/delegate-agent/.env
Environment=DATA_DIR=/opt/delegate-agent/data
ExecStart=/usr/bin/docker compose up
ExecStop=/usr/bin/docker compose down
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=token-dashboard

[Install]
WantedBy=multi-user.target
